# deploy/app-date.cron
2 2 * * * root /usr/sbin/timesync 2>/dev/null
